# ldm/crontab
MAILTO=""
# run ldm scour
0 0 * * * bin/ldmadmin scour
